# lib/systemd/system/postgresql@.service.d/foxxy_postgresql.conf
[Service]
ExecStartPre=/usr/share/foxxy-defaults/postgresql_reduce_shared_buffers %I
